# rc110_robot/rc110_system/systemd/rc110-roscore.service
[Unit]
Description=ROS Core

[Service]
Type=simple
ExecStart=/usr/lib/systemd/user/rc110-roscore.sh

[Install]
WantedBy=default.target
