# rc110_service/deb/postinst
#!/bin/bash

# Enable services by default
systemctl --global enable rc110
